# supervisord-conf/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/home/netalase/log/nginx.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=5
stderr_logfile=/home/netalase/log/nginx.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=5
